# deploy/cron/com.djinnfoundry.postiz-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd plist for macOS — postiz-agent dispatch daily at 08:00.

  Install:
    cp deploy/cron/com.djinnfoundry.postiz-agent.plist ~/Library/LaunchAgents/
    # edit the plist: replace /path/to/postiz-agent with your absolute path
    launchctl load ~/Library/LaunchAgents/com.djinnfoundry.postiz-agent.plist

  Verify:
    launchctl list | grep postiz-agent
    tail -f /path/to/postiz-agent/data/cron.log

  Unload:
    launchctl unload ~/Library/LaunchAgents/com.djinnfoundry.postiz-agent.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.djinnfoundry.postiz-agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/pnpm</string>
    <string>dev</string>
    <string>dispatch</string>
    <string>--platforms</string>
    <string>x,tiktok,instagram,youtube</string>
    <string>--json</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/path/to/postiz-agent</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/path/to/postiz-agent/data/cron.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/postiz-agent/data/cron.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
